# tools/mcp-governor/systemd/mcp-governor-report.timer
[Unit]
Description=Generate the MCP observation report daily

[Timer]
OnCalendar=*-*-* 00:15:00
Persistent=true
Unit=mcp-governor-report.service

[Install]
WantedBy=timers.target
